# configs/rock2/includes.chroot/etc/inittab
T0:2345:respawn:/sbin/getty -L ttyFIQ0 115200 vt100
