# config/includes.chroot/etc/rc.local
echo 1 > /sys/class/rkwifi/driver
ifconfig wlan0 up
